# ram-root/rc.local
#!/bin/sh

wifi down
/etc/init.d/wpad restart
wifi up

exit 0
